# 08_production/resources/etc/cron.d/simplesamlphp
# Run cron [daily]
02 0 * * * root curl --silent "https://sso-dev.fau.de/devssp/module.php/cron/cron.php?key=r34uef8ti6ti38396w7vt0nfcq62tgbb&tag=daily" > /dev/null 2>&1
# Run cron [hourly]
01 * * * * root curl --silent "https://sso-dev.fau.de/devssp/module.php/cron/cron.php?key=r34uef8ti6ti38396w7vt0nfcq62tgbb&tag=hourly" > /dev/null 2>&1