# 06_metarefresh/resources/etc/cron.d/simplesamlphp
# Run cron [daily]
02 0 * * * root curl -sS "https://sso-dev.fau.de/devssp/module.php/cron/run/daily/r34uef8ti6ti38396w7vt0nfcq62tgbb"
# Run cron [hourly]
01 * * * * root curl -sS "https://sso-dev.fau.de/devssp/module.php/cron/run/hourly/r34uef8ti6ti38396w7vt0nfcq62tgbb"
